# intelmap.service
[Unit]
Description=IntelMap Backend
After=network.target

[Service]
Type=simple
User=intelmap
Group=intelmap
WorkingDirectory=/opt/intelmap/backend
EnvironmentFile=/opt/intelmap/.env
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
